# build/services/GearShift/finish
#!/bin/bash -l

c_info "GearShift terminated"
