# build/services/GearShift/finish
#!/bin/bash

c_info "GearShift terminated"
