# debian-buster/build/rootfs/etc/gearbox/services/sshd/run
#!/bin/bash

HOST_KEYS="$(find /etc/ssh/ -name 'ssh_host*')"
if [ "${HOST_KEYS}" == "" ]
then
	echo "Gearbox: Generating missing SSH host keys..."
	/usr/bin/ssh-keygen -A
fi

echo "Gearbox: Starting SSHD init"
exec /usr/sbin/sshd -D -e

