# build/services/sshd/finish
#!/bin/bash -l

c_info "SSHD terminated"
