# alpine-3.5/build/rootfs/etc/gearbox/services/GearShift/finish
#!/bin/bash
echo "GearShift terminated"
